# steam-headless/overlay/etc/supervisor.d/gstreamer.ini
[program:gstreamer]
priority=30
autostart=true
autorestart=true
directory=/gst-plugins-rs
user=%(ENV_USER)s
command=GST_DEBUG=3 GST_PLUGIN_PATH="target/x86_64-unknown-linux-gnu/debug:$GST_PLUGIN_PATH" DISPLAY=:55 gst-launch-1.0     ximagesrc ! video/x-raw ! videoconvert ! capsfilter caps="video/x-raw,format=I420" ! nvh264enc ! h264parse config-interval=-1 ! queue max-size-buffers=500 ! livekit.video_0     pulsesrc ! audioconvert ! opusenc ! queue ! livekit.audio_0 livekitwebrtcsink name=livekit signaller::ws-url="REMOVED<" signaller::room-name="steam-test" signaller::api-key="%(ENV_LK_API_KEY)" signaller::secret-key="%(ENV_LK_SECRET_KEY)"
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
stopsignal=INT
stdout_logfile=/var/log/supervisor/gstreamer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/gstreamer.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7