# steam-headless/overlay/etc/supervisor.d/input-handler.ini
[program:input-handler]
priority=10
autostart=true
autorestart=true
directory=/input-handler
user=%(ENV_USER)s
command= go run main.go
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
stopsignal=INT
stdout_logfile=/var/log/supervisor/input-handler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/input-handler.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7